# postfix/etc/supervisord.d/postfix.ini
[program:postfix]
command=/usr/local/bin/postfix.sh
priority=100
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true