# postfix/etc/supervisord.d/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=999
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true